# tv-guide/services/tv-guide/run
#!/usr/bin/with-contenv bashio

# The app code is baked into the image at build time (see Dockerfile). Updating
# the add-on rebuilds the image from the repo, so code and version move together.

# ── Read config and export env vars ───────────────────────────────────────
export SONARR_URL="$(bashio::config 'sonarr_url')"
export SONARR_API_KEY="$(bashio::config 'sonarr_api_key')"
export TMDB_API_KEY="$(bashio::config 'tmdb_api_key')"
export HA_URL="$(bashio::config 'ha_url')"
export HA_TOKEN="$(bashio::config 'ha_token')"
export FIRETV_ENTITY="$(bashio::config 'firetv_entity')"
export SONOS_ENTITY="$(bashio::config 'sonos_entity')"

bashio::log.info "Sonarr: ${SONARR_URL} (key set: $([ -n "$SONARR_API_KEY" ] && echo yes || echo NO))"
bashio::log.info "TMDB key set: $([ -n "$TMDB_API_KEY" ] && echo yes || echo NO)"

bashio::log.info "Starting TV Guide on port 8099..."
exec uvicorn main:app --host 0.0.0.0 --port 8099 --app-dir /app
